# update-motd.d/10-hostname-color
#!/bin/bash

/usr/bin/env figlet "sisyphus@$(hostname)" | /usr/bin/env lolcat -f
